# cron/daily_journal.cron
# 每天23:50生成日记并同步到GitHub
50 23 * * * cd /home/ubuntu/.openclaw/workspace && python3 scripts/daily_journal.py >> logs/journal.log 2>&1
51 23 * * * cd /home/ubuntu/.openclaw/workspace && bash scripts/sync_to_github.sh >> logs/github_sync.log 2>&1